# crontab
*/5 * * * * cd /app/ && python3 -m scraper.main --dst-email="omstrumpf_7956f1@kindle.com" > /proc/1/fd/1 2> /proc/1/fd/2
